# udev.rules.d/99-backlight.rules
SUBSYSTEM=="backlight", ACTION=="add", RUN+="/usr/bin/chmod o+w /sys/class/backlight/%k/brightness"
